# monitor.service
[Unit]
# openai_register 持续巡检 / 自动补号服务
Description=openai_register monitor service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/data/openai_register
# 如需代理、阈值、自定义配置文件，请修改 monitor_config.json
# 可追加 --dingtalk-webhook <URL> 设置钉钉通知
ExecStart=/data/openai_register/.venv/bin/python /data/openai_register/openai_register.py \
    --config /data/openai_register/monitor_config.json
Restart=always
RestartSec=10
# 日志由 journald 管理，自动轮转，无需手动处理日志文件
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openai-register
Environment=LANG=C
Environment=LC_ALL=C
# 所有配置项（含钉钉 webhook）均通过 monitor_config.json 管理

[Install]
WantedBy=multi-user.target
