# monitor.service
[Unit]
# openai_register 持续巡检 / 自动补号服务
Description=openai_register monitor service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/data/openai_register
# 请先执行 bash ctl.sh deps 或 bash ctl.sh install 创建 .venv
# 当前服务为 A-only 巡检模式；运行参数统一放在 monitor_config.json
# 不要再追加旧版参数：--pool-min-count / --usage-threshold
ExecStart=/data/openai_register/.venv/bin/python /data/openai_register/openai_register.py --monitor --config /data/openai_register/monitor_config.json
Restart=always
RestartSec=10
# 日志由 journald 管理，自动轮转，无需手动处理日志文件
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openai-register
Environment=LANG=C
Environment=LC_ALL=C
Environment=PYTHONUNBUFFERED=1
# 所有配置项（含钉钉 webhook）均通过 monitor_config.json 管理

[Install]
WantedBy=multi-user.target
